# PiController/adc-daemon/adc-daemon.service
[Unit]
Description=ADC Daemon Service
After=network.target

[Service]
Type=simple
ExecStart=/opt/garden/python3/bin/python3 /opt/garden/bin/adc-daemon.py --hwversion 5.1
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
